# sysstatsapp/etc/systemd/system/sysstatsapp.service
[Unit]
Description=CPU and RAM Usage Logger

[Service]
Type=simple
ExecStart=/bin/sysstatsapp.sh

[Install]
WantedBy=multi-user.target